# root/etc/services.d/ttrss-updater/run
#!/usr/bin/with-contenv sh
sleep 24h
exec s6-setuidgid www-data /srv/update-ttrss.sh
# This will automatically restart, wait 24h again and reupdate
